# etc/s6-overlay/s6-rc.d/reproxy-log-prepare/up
if { mkdir -p /var/log/reproxy }
if { chown nobody:nogroup /var/log/reproxy }
chmod 02755 /var/log/reproxy
