# etc/s6-overlay/s6-rc.d/reproxy/run
#!/command/execlineb -P
reproxy \
    --logger.enabled \
    --logger.file=/var/log/reproxy/access.log \
    --assets.location=/var/www/galaiko.rocks \
    --listen=0.0.0.0:8080
