# etc/s6-overlay/s6-rc.d/goaccess/run
#!/command/execlineb -P
goaccess --ignore-referrer=*.galaiko.rocks
